# dmarc-analyzer/crontab
0 0 * * * /usr/local/bin/dmarc-report-converter -config /config.yaml